# nodes/rt-radio/systemd/rt-radio-presence.service
[Unit]
Description=RollingThunder rt-radio presence publisher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=spiff
Group=spiff
ExecStart=/usr/bin/python3 /opt/rollingthunder/nodes/rt-radio/services/rt-radio_node_presence_publisher.py
Restart=always
RestartSec=2
Environment=RT_MQTT_HOST=rt-controller
Environment=RT_MQTT_PORT=1883
Environment=RT_PRESENCE_HTTP_PORT=8787
Environment=RT_NODE_ID=rt-radio
Environment=RT_NODE_ROLE=radio

[Install]
WantedBy=multi-user.target

