# nodes/rt-radio/systemd/rt-radio-ui-intent-worker.service
[Unit]
Description=RollingThunder UI Intent Worker (rt-radio)
After=network-online.target redis-server.service rigctld.service
Wants=network-online.target rigctld.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/rollingthunder

Environment=RT_NODE_ID=rt-radio
Environment=RT_UI_INTENTS_CHANNEL=rt:ui:intents
Environment=RT_UI_BUS_CHANNEL=rt:ui:bus
Environment=RT_CONFIG_PATH=/opt/rollingthunder/config/app.json
Environment=RT_REDIS_HOST=rt-controller
Environment=RT_REDIS_PORT=6379
Environment=RT_REDIS_DB=0
Environment=RT_ALLOW_REBOOT=1
EnvironmentFile=/etc/rollingthunder/redis.env
EnvironmentFile=-/etc/rollingthunder/radio.env

ExecStart=/opt/rollingthunder/.venv/bin/python /opt/rollingthunder/tools/ui_intent_worker.py

# Hardening tuned for current behavior
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target