# nodes/rt-controller/systemd/rt-alert@.service
[Unit]
Description=RollingThunder Alert Writer for %I

[Service]
Type=oneshot
User=spiff
ExecStart=/usr/bin/python3 /opt/rollingthunder/services/rt_emit_alert.py --source systemd --service "%I" --level bad --message "Service failed and was stopped after restart limits"