# nodes/rt-controller/systemd/rt-alert@.service
[Unit]
Description=RollingThunder Alert Writer for %I

[Service]
Type=oneshot
User=spiff
WorkingDirectory=/opt/rollingthunder
ExecStart=/usr/bin/python3 /opt/rollingthunder/services/rt_emit_alert.py \
  --source systemd \
  --kind service \
  --severity bad \
  --service "%I" \
  --title "Service failed" \
  --message "Service failed and was stopped after restart limits"
  