# nodes/rt-controller/systemd/rt-controller-ui-intent-worker.service
[Unit]
Description=RollingThunder UI Intent Worker (rt-controller)
After=network-online.target redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=spiff
Group=spiff
WorkingDirectory=/opt/rollingthunder

Environment=RT_NODE_ID=rt-controller
Environment=RT_UI_INTENTS_CHANNEL=rt:ui:intents
Environment=RT_UI_BUS_CHANNEL=rt:ui:bus
Environment=RT_CONFIG_PATH=/opt/rollingthunder/config/app.json
EnvironmentFile=/etc/rollingthunder/redis.env

ExecStart=/usr/bin/python3 /opt/rollingthunder/services/ui_intent_worker.py

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target