# nodes/rt-display/systemd/rt-display-presence.service
[Unit]
Description=RollingThunder rt-display presence publisher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=spiff
Group=spiff
ExecStart=/opt/rollingthunder/.venv/bin/python /opt/rollingthunder/nodes/rt-display/services/display_presence.py
Restart=always
RestartSec=2
Environment=RT_MQTT_HOST=rt-controller
Environment=RT_MQTT_PORT=1883
Environment=RT_PRESENCE_HTTP_PORT=8787


[Install]
WantedBy=multi-user.target
