# nodes/rt-display/ops/rt-display-ui.service.template
[Unit]
Description=RollingThunder rt-display UI static server
After=network.target

[Service]
Type=simple
User=spiff
WorkingDirectory=/opt/rollingthunder/nodes/rt-display/ui
ExecStart=/usr/bin/python3 -m http.server 8619 --bind 0.0.0.0
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
