# nodes/rt-wpsd/systemd/rt-presence-publisher.service
[Unit]
Description=RollingThunder Presence Publisher (MQTT, read-only)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi-star
Environment=RT_NODE_ID=ki5vnb-dmr2
Environment=RT_NODE_ROLE=external
Environment=RT_MQTT_HOST=rt-controller
Environment=RT_MQTT_PORT=1883
Environment=RT_PRESENCE_TOPIC_PREFIX=rt/presence
EnvironmentFile=/etc/rollingthunder/redis.env
ExecStart=/opt/rollingthunder/tools/publish_presence.sh
