# deploy/nodes/rt-controller/systemd/rollingthunder-controller.service
[Unit]
Description=RollingThunder Controller Bootstrap (rt-controller)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

User=spiff
Group=spiff
WorkingDirectory=/opt/rollingthunder

ExecStart=/opt/rollingthunder/.venv/bin/python /opt/rollingthunder/nodes/rt-controller/rt_controller.py

Restart=on-failure
RestartSec=2

Environment=PYTHONUNBUFFERED=1
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
