# deploy/nodes/rt-controller/systemd/rt-env-temp-publisher.service
[Unit]
Description=RollingThunder Env Temp Publisher
After=network.target redis-server.service
Wants=network.target

[Service]
Type=simple
User=root

# 🔑 Load Redis credentials (same pattern as other RT services)
EnvironmentFile=/etc/rollingthunder/redis.env

ExecStart=/opt/rollingthunder/.venv/bin/python /opt/rollingthunder/services/env_temp_publisher.py

Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target