# deploy/nodes/rt-controller/systemd/rt-pota-context-manager.service
[Unit]
Description=RollingThunder POTA Context Manager
After=network-online.target redis.service
Wants=network-online.target

[Service]
Type=simple
User=spiff
Group=spiff
WorkingDirectory=/opt/rollingthunder
ExecStart=/opt/rollingthunder/.venv/bin/python /opt/rollingthunder/services/pota_context_manager.py
Restart=on-failure
RestartSec=2

Environment=RT_REDIS_HOST=127.0.0.1
Environment=RT_REDIS_PORT=6379
Environment=RT_REDIS_DB=0
EnvironmentFile=-/etc/rollingthunder/redis.env

Environment=RT_POTA_CONTEXT_KEY=rt:pota:context
Environment=RT_POTA_NEARBY_KEY=rt:pota:nearby
Environment=RT_POTA_SSB_BANDS_SOURCE_KEY=rt:pota:ssb:bands
Environment=RT_POTA_SSB_SPOTS_SOURCE_PREFIX=rt:pota:ssb:spots
Environment=RT_POTA_UI_BANDS_KEY=rt:pota:ui:ssb:bands
Environment=RT_POTA_UI_SPOTS_PREFIX=rt:pota:ui:ssb:spots

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target