# deploy/nodes/rt-controller/systemd/rt-wpsd-log-ingestor.service
[Unit]
Description=RollingThunder WPSD Log Ingestor (read-only WPSD -> Redis)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=spiff
Environment=PYTHONUNBUFFERED=1
Environment=RT_WPSD_HOST=192.168.8.184
Environment=RT_WPSD_USER=pi-star
# Optional hardening if you add these env hooks in the script:
# Environment=RT_WPSD_SSH_KEY=/home/spiff/.ssh/id_ed25519
# Environment=RT_WPSD_KNOWN_HOSTS=/home/spiff/.ssh/known_hosts

ExecStart=/opt/rollingthunder/.venv/bin/python /opt/rollingthunder/services/wpsd_log_ingestor.py
Restart=on-failure
RestartSec=2
TimeoutStartSec=10

[Install]
WantedBy=multi-user.target
