# deploy/nodes/rt-controller/systemd/rt-wpsd-poller.service
[Unit]
Description=RollingThunder WPSD Poller (read-only WPSD -> Redis)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=spiff
Environment=PYTHONUNBUFFERED=1
Environment=RT_WPSD_HOST=192.168.8.184
Environment=RT_WPSD_USER=pi-star

ExecStart=/opt/rollingthunder/.venv/bin/python /opt/rollingthunder/services/wpsd_poller.py
Restart=on-failure
RestartSec=2
TimeoutStartSec=10

[Install]
WantedBy=multi-user.target
