# nodes/rt-display/autostart_systemd/rt-display-kiosk-watchdog.service
[Unit]
Description=RollingThunder rt-display Kiosk Watchdog
After=rt-display-kiosk.service
Requires=rt-display-kiosk.service

[Service]
Type=simple
WorkingDirectory=/opt/rollingthunder
ExecStart=/opt/rollingthunder/nodes/rt-display/services/rt-display-kiosk-watchdog.sh

Restart=always
RestartSec=15

Environment=RT_WATCHDOG_MEM_MB=1100
Environment=RT_WATCHDOG_SWAP_MB=350
Environment=RT_WATCHDOG_CPU_PCT=95
Environment=RT_WATCHDOG_CPU_SEC=180
Environment=RT_WATCHDOG_INTERVAL=60

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target