# nodes/rt-display/autostart_systemd/rt-display-kiosk.service
[Unit]
Description=RollingThunder rt-display Wayland Chromium Kiosk
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/opt/rollingthunder/nodes/rt-display/ops/rt-display-kiosk.sh
Restart=always
RestartSec=5
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=WAYLAND_DISPLAY=wayland-0

[Install]
WantedBy=default.target